# recipes-core/updatehub/updatehub/updatehub.service
[Unit]
Description=updatehub agent
After=local-fs.target network.target time-sync.target
Requires=local-fs.target

[Service]
ExecStart=@BINDIR@/unshare -m @BINDIR@/updatehub
Restart=on-failure

[Install]
WantedBy=multi-user.target
